# overlay/etc/s6-overlay/s6-rc.d/node/run
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Run preflight checks
#======================================================================================================================

nu -c "use bf-node run ; run preflight"


#======================================================================================================================
# Start Node.js server
#======================================================================================================================

s6-setuidgid node node ${BF_NODE_RUN}
